# bionic/edit/etc/systemd/system/cybereason-sensor.service.d/override.conf
[Service]
Environment="http_proxy=http://dyn.wmes.pl:8080"
Environment="https_proxy=http://dyn.wmes.pl:8080"
